# deploy/cc-core.service
# CC backend (FastAPI / uvicorn) — Story 1.7
#
# Install: sudo cp deploy/cc-core.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now cc-core
#
# SINGLE-OWNER RULE (architecture mandate): from Story 2.2 on, this process
# owns the Telegram session file (/var/lib/cc/anon.session). NEVER run a
# second instance of this service — two MTProto clients on one session file
# corrupt it and can invalidate the Telegram auth.
#
# backend/.env is loaded by app/config.py itself (path resolved relative to
# the module, CWD-independent) — no EnvironmentFile= needed here.

[Unit]
Description=CC backend (FastAPI/uvicorn)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=cc
Group=cc
WorkingDirectory=/srv/cc/backend
# Bind 127.0.0.1 ONLY — Caddy is the sole public listener. No --reload in prod.
ExecStart=/srv/cc/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
